# install/systemd.target
[Unit]
Wants=dokku-event-listener.service

[Install]
WantedBy=multi-user.target
